# docker-monitoring/bin/monitoring.daemon.service.template
[Unit]
Description=cAdvisor and Node Export Monitoring Service
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{RunnerRoot}}
ExecStart={{Command}} up -d
ExecStop={{Command}} down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target